# server-config/baxter_flask.conf
[program:baxter_flask]
command = /home/www/gunicorn-start-baxter
stdout_logfile = /home/www/logs/baxter-flask/gunicorn_supervisor.log
redirect_stderr = true
